# files/uzblmonitor.conf
description "Start uzbl yo"
author "Patrick Lucas <me@patricklucas.com>"

start on desktop-session-start

env DISPLAY=:0.0

script
    exec sudo -u monitor -i /usr/bin/uzblmonitor
end script
